# board/licheepi/rootfs-overlay/etc/profile.d/env_vars.sh
#export LC_ALL=en_US.UTF-8
export LANG=ru_RU.UTF-8

setfont Cyr_a8x14.psfu.gz

loadkeys ruwin_ct_sh-UTF-8

